# rc.d/dbstore.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: dbstore
# REQUIRE: LOGIN

. /etc/rc.subr

name="dbstore"
rcvar="dbstore_enable"

pidfile="@app_rundir@"/dbstore.pid
command="@prefix@/sbin/${name}"
command_args=""
procname="@prefix@/sbin/${name}"
#command_interpreter="@perl@"

load_rc_config ${name}

: ${dbstore_enable:="NO"}

run_rc_command "$1"
#EOF
